# PMEngined-0.1.0/starter
#!/sbin/runscript

start() {
	ebegin "Starting watchdog"
	start-stop-daemon --start --quiet --exec PMEngined test.ini \

	eend ${?}
}

stop() {
	ebegin "Stopping watchdog"
	start-stop-daemon --stop --retry 5 --quiet --exec PMEngined
	eend ${?}
}
